# root/etc/services.d/mosdns-watch/run
#!/usr/bin/with-contenv sh

while inotifywait -e modify --exclude '.*\.log$' /config; do
    s6-svc -r /var/run/s6/services/mosdns
done