# root/etc/services.d/mosdns/run
#!/usr/bin/with-contenv sh

exec /usr/bin/mosdns -dir /config